# setup/kibana.service
[Unit]
Description=Kibana
After=network.target

[Service]
ExecStart=/home/vagrant/kibana-SNAPSHOT/bin/kibana
Type=simple
User=vagrant
PIDFile=/run/kibana.pid
Restart=always

[Install]
WantedBy=default.target
